# backend/crontab
*/30 * * * * /usr/local/bin/python3 /cron-jobs/manage.py ats_jobs > /proc/1/fd/1 2>/proc/1/fd/2
